# tool/mesos/docker/singa/ssh.conf
[supervisord]
nodaemon=true

[program:ssh]
command=/usr/sbin/sshd -D
